# units/etcd.service
[Unit]
Description=etcd
Documentation=https://github.com/coreos

[Service]
EnvironmentFile=-/etc/etcd.env
ExecStart=/usr/local/bin/etcd \
  --name master3 \
  --cert-file=/etc/etcd/kube-api-server.crt \
  --key-file=/etc/etcd/kube-api-server.key \
  --peer-cert-file=/etc/etcd/kube-api-server.crt \
  --peer-key-file=/etc/etcd/kube-api-server.key \
  --trusted-ca-file=/etc/etcd/ca.crt \
  --peer-trusted-ca-file=/etc/etcd/ca.crt \
  --initial-advertise-peer-urls https://192.168.100.151:2380 \
  --listen-peer-urls https://192.168.100.151:2380 \
  --listen-client-urls https://192.168.100.151:2379,https://127.0.0.1:2379 \
  --advertise-client-urls https://192.168.100.151:2379 \
  --initial-cluster-token etcd-cluster-0 \
  --initial-cluster master1=https://192.168.100.149:2380,master2=https://192.168.100.150:2380,master3=https://192.168.100.151:2380 \
  --initial-cluster-state new \
  --data-dir=/var/lib/etcd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
